# pypi-server.conf
[program:pypi-server]
directory=/home/pypi/
command=/home/pypi/venv/bin/pypi-server -p 10086 -P /home/pypi/.htaccess /home/pypi/packages
autostart=true
autorestart=true
redirect_stderr=true
